# synexs-ai-swarm.service
[Unit]
Description=Synexs AI Swarm
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/synexs
Environment="WORK_DIR=/root/synexs"
Environment="CYCLE_INTERVAL=7200"
ExecStart=/root/synexs/synexs_env/bin/python3 -u /root/synexs/ai_swarm_fixed.py
Restart=on-failure
RestartSec=300
StandardOutput=append:/var/log/synexs-ai-swarm.log
StandardError=append:/var/log/synexs-ai-swarm.log

[Install]
WantedBy=multi-user.target
