# synexs-honeypot.service
[Unit]
Description=Synexs Honeypot Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/synexs
ExecStart=/usr/bin/python3 /root/synexs/honeypot_server.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/synexs-honeypot.log
StandardError=append:/var/log/synexs-honeypot.log

[Install]
WantedBy=multi-user.target
